# zeppelin-spark-centos/roles/zeppelin/templates/zeppelin.service
[Unit]
Description=Zeppelin Service

[Service]
ExecStart=/opt/zeppelin/bin/zeppelin.sh

[Install]
WantedBy=multi-user.target
